# docker_files/apache2/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/apache2 -D FOREGROUND

[program:sshd]
command=/usr/sbin/sshd -D