# docker_files/memcached/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:memcached]
command=/usr/bin/memcached -u root -m 64
